# utilities/piFakeBoot/measure.eltt2.service
[Unit]
Description=Fake Measured Boot Eltt2 version
#Wants=tpm2-abrmd.service
#After=tpm2-abrmd.service

[Service]
Type=oneshot
ExecStart=/boot/measuredboot/measure.start.eltt2

[Install]
WantedBy=multi-user.target
