# dist/ga10.service
[Unit]
Description=GA10 Attestation Engine
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=ian
ExecStart=/opt/nae/ga10 -config=/opt/nae/ga10config.yaml

[Install]
WantedBy=multi-user.target